# infra/systemd/hantaatlas-worker.service
[Unit]
Description=HantaAtlas official-source ingestion worker
After=network-online.target postgresql.service hantaatlas-api.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=hantaatlas
Group=hantaatlas
WorkingDirectory=/opt/hantaatlas
ExecStart=/usr/bin/node /opt/hantaatlas/backend/worker/dist/index.js
EnvironmentFile=/etc/hantaatlas/secrets.env

Restart=on-failure
RestartSec=15s
StartLimitIntervalSec=120
StartLimitBurst=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
PrivateUsers=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0027

ReadWritePaths=/var/log/hantaatlas
ReadOnlyPaths=/etc/hantaatlas

LimitNOFILE=65535
TasksMax=256

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hantaatlas-worker

[Install]
WantedBy=multi-user.target
