# infra/systemd/hantaatlas-worker.timer
[Unit]
Description=Run HantaAtlas worker every hour

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h
RandomizedDelaySec=2min
Persistent=true

[Install]
WantedBy=timers.target
