# infra/systemd/hantaatlas-backup.timer
[Unit]
Description=Run hantaatlas-backup nightly

[Timer]
OnCalendar=daily
RandomizedDelaySec=15min
Persistent=true

[Install]
WantedBy=timers.target
